# deploy/systemd/poly_5min_bot.service
[Unit]
Description=Polymarket 5-minute arbitrage bot
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=polybot
Group=polybot
WorkingDirectory=/opt/poly_5min_bot
EnvironmentFile=/etc/poly_5min_bot/poly_5min_bot.env
ExecStart=/opt/poly_5min_bot/target/release/poly_5min_bot
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGINT

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
UMask=0077
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
